# configs/lastwatt.service
[Unit]
Description=LastWatt Grid Curtailment Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/mcd/lastwatt/lastwatt daemon --config /home/mcd/lastwatt/config.yaml
Restart=always
# Don't restart on EX_CONFIG (78) — a bad config won't fix itself by retrying;
# the daemon exits 78 on config errors so the unit fails cleanly instead of looping.
RestartPreventExitStatus=78
RestartSec=10
User=mcd
Group=mcd

# Allow unprivileged ICMP ping
AmbientCapabilities=CAP_NET_RAW

# State directory
StateDirectory=lastwatt

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/lastwatt

[Install]
WantedBy=multi-user.target
